# config/supervisor.ini
[program:kool967]
command=env/bin/gunicorn server:app -c config/gunicorn.conf.py
directory=/home/bigtree/dev/python/kool967
user=bigtree
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=log/access.log
stderr_logfile=log/error.log
